# crontab
11 11 ** 4  /home/$USER/.kosmosautomat/updateScript.sh

#####################################
# Schedule at 11:11 am every Thursday 
# Of Each Week .
#####################################